# minecraft.service
# Servicio systemd para Minecraft (PaperMC) con consola en screen
# Autor: Jota (Jota-FalseProphet)
# Uso:
#   - Copia este archivo como /etc/systemd/system/minecraft.service
#   - systemctl daemon-reload
#   - systemctl enable --now minecraft
#   - Para ver la consola: sudo screen -r minecraft
#   - Para salir sin parar el server: Ctrl+A, luego D

[Unit]
Description=Minecraft Paper Server (con screen)
After=network.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/minecraft/paper
User=yotesis
Group=yotesis

Environment="JAVA_ARGS=-Xms4G -Xmx6G -XX:+UseG1GC -XX:MaxGCPauseMillis=100 -XX:+ParallelRefProcEnabled -XX:+UnlockExperimentalVMOptions -XX:+DisableExplicitGC -XX:+AlwaysPreTouch -XX:G1NewSizePercent=30 -XX:G1MaxNewSizePercent=40 -XX:G1HeapRegionSize=8M -XX:G1ReservePercent=20 -XX:G1HeapWastePercent=5 -XX:G1MixedGCCountTarget=4 -XX:InitiatingHeapOccupancyPercent=15 -XX:G1MixedGCLiveThresholdPercent=90 -XX:G1RSetUpdatingPauseTimePercent=5 -XX:SurvivorRatio=32 -XX:+PerfDisableSharedMem -XX:MaxTenuringThreshold=1"
Environment="SCREENDIR=/opt/minecraft/.screen"

ExecStartPre=/bin/mkdir -p /opt/minecraft/.screen
ExecStartPre=/bin/chown -R yotesis:yotesis /opt/minecraft/.screen
ExecStart=/usr/bin/screen -DmS minecraft /usr/bin/java $JAVA_ARGS -jar server.jar --nogui
ExecStop=/usr/bin/screen -S minecraft -p 0 -X stuff "stop$(printf 
)"
ExecStop=/bin/sleep 2
ExecStop=/usr/bin/screen -S minecraft -X quit

Restart=on-failure
RestartSec=10
SuccessExitStatus=0 143
Nice=5

[Install]
WantedBy=multi-user.target
